# .docker/crontab/cronjob.alpine
1 0 * * * cd /var/www/kanboard && ./kanboard cronjob >/dev/null 2>&1
