# docker/crontab/cronjob.alpine
1 0 * * * cd /var/www/app && ./kanboard cronjob
