# .docker/services.d/nginx/run
#!/bin/execlineb -P
nginx -g "daemon off;"