# docker/etc/crontabs/nginx
0 8 * * * cd /var/www/app && php84 ./cli cronjob >/dev/null 2>&1
